# templates/default/systemd.service.erb
[Unit]
Description=<%= @desc -%>
After=syslog.target

[Service]
User=<%= @user -%>
ExecStart=<%= @exec_start -%>
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
